# extra/supervisor.conf
[program:pf]
user = root
group = root
command = python test.py
directory = /var/lib/container/pf
autostart = true
autorestart = true
stopsignal = 15
pidfile = /var/run/container/pf/test.pid
